# mystvpn.service
[Unit]
Description=MystVPN Telegram Bot
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/MystBot
ExecStart=/home/ubuntu/MystBot/.venv/bin/python main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Автоматически перезапускаться при краше
StartLimitIntervalSec=60
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
